# watchdog-restart.timer
Unit]
Description=Restart USB watchdog each minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
